# configs/com.almacreativa.hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.almacreativa.hermes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/hermes-start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{HOME}}/.hermes</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HERMES_HOME</key>
        <string>{{HOME}}/.hermes</string>
        <key>HERMES_DASHBOARD_HOST</key>
        <string>0.0.0.0</string>
        <key>HERMES_DASHBOARD_PORT</key>
        <string>9119</string>
        <key>PATH</key>
        <string>{{HOME}}/.local/bin:{{HOME}}/.opencode/bin:{{HOME}}/.nvm/versions/node/{{NODE_VERSION}}/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/ai-lab/logs/hermes.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/ai-lab/logs/hermes.log</string>
</dict>
</plist>
