# configs/dagu.service
[Unit]
Description=Dagu Workflow Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{DAGU_PATH}} start-all --config={{HOME}}/.config/dagu/config.yaml
Restart=always
RestartSec=10
Environment=HOME={{HOME}}

[Install]
WantedBy=default.target
